# ha-addon/calendhd/rootfs/etc/services.d/pocketbase/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# calenDHD PocketBase service
# ==============================================================================
exec /run.sh
